# scripts/init/dime-server.conf
description "Digital Me server"

start on runlevel [2345]
stop on runlevel [06]

exec $DIME_INSTALL_DIR/bin/run-dime.sh
